# deploy/schoolair-dev.service.example
[Unit]
Description=SchoolAir sensor ingest (DEV)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/repo_path/pi
ExecStart=%h/repo_path/pi/.venv/bin/python -u main.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target