# deploy/schoolair.service
[Unit]
Description=SchoolAir sensor ingest
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/schoolair
ExecStart=/home/pi/schoolair/.venv/bin/python -u /home/pi/schoolair/main.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
