# deploy/schoolair-netwatch.service
[Unit]
Description=SchoolAir Network Watchdog
Documentation=https://github.com/SchoolAir/schoolair-ex-RMIT-pi
# Start after the boot launcher has finished its one-shot decision.
After=network.target NetworkManager.service schoolair-launcher.service

[Service]
Type=simple
ExecStart=/home/admin/schoolair/registration_wizard/netwatch.sh
User=root
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=schoolair-netwatch

[Install]
WantedBy=multi-user.target
