# deploy/schoolair-launcher.service
[Unit]
Description=SchoolAir Network Launcher
# Run after NetworkManager is up so nmcli queries work.
After=NetworkManager.service schoolair-first-boot.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/home/admin/schoolair/registration_wizard/launcher.sh
RemainAfterExit=no
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
