# deploy/schoolair.service
[Unit]
Description=SchoolAir sensor ingest
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/schoolair
ExecStart=/home/admin/schoolair/.venv/bin/python -u /home/admin/schoolair/main.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
