# deploy/sen6x.service
[Unit]
Description=Sensirion SEN6x Sensor — one-shot initialisation
# Runs sen6x_read --init once at boot: if the sensor is already in continuous
# measurement mode it returns immediately; otherwise it resets, starts
# continuous measurement, and waits up to 75 s for the first valid CO2 sample.
# schoolair.service depends on this completing successfully before starting.
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/admin/i2c/sen6x/sen6x_read --init
User=admin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
